# src/component/custom/etc/files/profile.d/recipe.sh
export RECIPE_PATH=/etc/recipe/repos/rocky10
