# profile/live-overlay/etc/systemd/system/getty@tty1.service.d/live-disable.conf
[Unit]
ConditionPathExists=!/live-boot-marker
